# gateway-core/deploy/matters-gateway-cloudwatch-metrics.service.example
[Unit]
Description=Publish Matters Fediverse gateway runtime metrics to CloudWatch
Requires=matters-gateway-core.service
After=matters-gateway-core.service network-online.target

[Service]
Type=oneshot
User=root
Environment=AWS_REGION=ap-southeast-1
Environment=CLOUDWATCH_NAMESPACE=Matters/FediverseGateway
WorkingDirectory=/opt/matters-gateway/repo/gateway-core
ExecStart=/usr/bin/env node scripts/publish-cloudwatch-metrics.mjs --config /etc/matters-gateway/staging.instance.json
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
