# gateway-core/deploy/matters-gateway-cloudwatch-metrics.timer.example
[Unit]
Description=Publish Matters Fediverse gateway runtime metrics every five minutes

[Timer]
OnBootSec=3m
OnUnitActiveSec=5m
AccuracySec=30s
RandomizedDelaySec=15s
Unit=matters-gateway-cloudwatch-metrics.service

[Install]
WantedBy=timers.target
