# gateway-core/deploy/matters-gateway-delivery.timer.example
[Unit]
Description=Retry pending Matters Fediverse deliveries every five minutes

[Timer]
OnBootSec=2m
OnUnitActiveSec=5m
AccuracySec=30s
RandomizedDelaySec=15s
Unit=matters-gateway-delivery.service

[Install]
WantedBy=timers.target
